# docker/supervisor/conf.d/nextjs.conf
[program:nextjs]
directory=/app/
command=npm start
startsecs = 0
autostart=true
autorestart=true
stderr_logfile=/var/log/perfectpoint/app.log
stdout_logfile=/var/log/perfectpoint/app.log